# systemd/homelab-gitpull.service
[Unit]
Description=Update homelab from git and redeploy
Wants=network-online.target
After=network-online.target

[Service]
User=groot
Group=docker
Type=oneshot
WorkingDirectory=/srv/homelab
ExecStart=/srv/homelab/scripts/deploy.sh
StandardOutput=journal
